# TemporaryAdmin/pkg_scripts/postinstall
#!/bin/bash

/bin/launchctl load /Library/LaunchDaemons/com.github.ygini.temporaryadmin.plist

exit 0
