# EnrollMe/pkg_scripts/postinstall
#!/bin/bash

/bin/launchctl load /Library/LaunchDaemons/com.github.ygini.enrollme.plist

exit 0
